# kind/hybrid-node-setup.service
[Unit]
Description=Hybrid Node Setup Service
After=network.target containerd.service
Wants=containerd.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/hybrid-node-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target